# roxy.service
[Unit]
Description=Roxy the Frontend Proxy
After=network.target

[Service]
Type=notify
NotifyAccess=main
EnvironmentFile=/etc/default/roxy
ExecStart=/opt/roxy/bin/roxy $ROXY_FLAGS
ExecReload=/opt/roxy/bin/roxyctl $ROXYCTL_FLAGS reload
ExecStop=/opt/roxy/bin/roxyctl $ROXYCTL_FLAGS shutdown
Restart=always
User=roxy
Group=roxy
UMask=022

[Install]
WantedBy=multi-user.target

; vim:set filetype=systemd:
